# shingo-core/deploy/shingo-core.service
[Unit]
Description=ShinGo Core
After=network-online.target
Wants=network-online.target
# Never give up restarting. A plant-floor service that quietly stays
# down is worse than a noisy crashloop — the latter is caught by
# journal/monitoring. Setting the interval to 0 disables the rate
# limit; StartLimitBurst is meaningless without it, so we drop it.
StartLimitIntervalSec=0

[Service]
Type=simple
User=shingo
Group=shingo
WorkingDirectory=/opt/shingo
ExecStart=/opt/shingo/shingocore --config /etc/shingo/shingocore.yaml
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shingo-core
MemoryMax=4G
LogsDirectory=shingo-core
Environment=SHINGO_PANIC_LOG=/var/log/shingo-core/last-panic.log
# Teams crash alert. Fires on every stop; alert-on-stop.sh handles its
# own suppression (deploy marker, system reboot, backoff window, empty
# webhook URL). Argument list is unit name + systemd-exported result vars.
ExecStopPost=/opt/shingo/alert-on-stop.sh %n ${SERVICE_RESULT} ${EXIT_CODE} ${EXIT_STATUS}

[Install]
WantedBy=multi-user.target
